# deploy/systemd/ipastore-scinsta-cancel@.service
[Unit]
Description=Abort a running SCInsta build for %i
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=altuser
Group=altuser
# Le flag de cancel est supprime par cmd_scinsta_cancel lui-meme (avant le
# docker kill), pour eviter que ce path unit retrigger apres le kill.
ExecStart=/usr/local/bin/website-management %i-scinsta-cancel
StandardOutput=journal
StandardError=journal
# Le kill lui-meme est rapide, mais docker peut prendre jusqu'a 10s pour
# propager SIGTERM puis SIGKILL. 30s laisse de la marge.
TimeoutStartSec=30
