# deploy/ipastore.service
[Unit]
Description=IPA Store - FastAPI admin panel
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=altuser
Group=altuser
WorkingDirectory=/opt/ipastore
Environment="PATH=/opt/ipastore/venv/bin:/usr/bin:/bin"
Environment="IPASTORE_STORE_DIR=/srv/store"
Environment="IPASTORE_SECRET_FILE=/opt/ipastore/secret_key"
ExecStart=/opt/ipastore/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/srv/store /opt/ipastore /var/log/ipastore

[Install]
WantedBy=multi-user.target
